# paco/crontab
*/5 * * * * rsync -a /srv/swift/ /etc/swift/
# Mandatory blank line
